# frontend/userpage/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn -w 4 -b :8000 main:app
directory=/home/ubuntu/StepUpSage_FE/userpage
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log
